# clamdscan/clamav-clamdscan@.timer.in
[Unit]
Description=Weekly ClamAV clamdscan %I

[Timer]
OnCalendar=weekly
AccuracySec=1h
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
