# freshclam/clamav-freshclam.service.in
[Unit]
Description=ClamAV virus database updater
Documentation=man:freshclam(1) man:freshclam.conf(5) https://docs.clamav.net/
# If user wants it run from cron, don't start the daemon.
ConditionPathExists=!/etc/cron.d/clamav-freshclam
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=@prefix@/bin/freshclam -d --foreground=true

##
## Security Hardening Options
##
ProtectSystem=full
ProtectHome=tmpfs
ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes

NoExecPaths=/
# If you want to run commands on event, append the full path of the command or
# executable to `ExecPaths`. Commonly, this is used for `OnUpdateExecute`,
# `OnErrorExecute`, and `OnOutdatedExecute` options in freshclam.conf. Make sure
# there is only one `ExecPaths` option. The binaries must be space separated like so:
;ExecPaths=@prefix@/sbin/freshclam /usr/local/bin/OnOutdated.sh /usr/local/bin/OnErrorExecute.sh
ExecPaths=@prefix@/bin/freshclam

[Install]
WantedBy=multi-user.target
